# app/supervisor.conf
; STT-Endpoint all-in-one container — supervisord program definitions.
;
; Order matters: stt first (FastAPI is needed by auth's startup readiness check),
; then auth (waits for STT /ready in start.sh), then nginx binds the public port.
;
; All programs log to stdout/stderr; the platform harvests them as one stream.

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/dev/null
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:stt]
command=/usr/local/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --log-level info
directory=/opt/stt/app
environment=
    STT_MODELS_DIR="/var/lib/stt/models",
    PYTHONUNBUFFERED="1",
    PYTHONDONTWRITEBYTECODE="1",
    ; Defaults from app/config.py are fine; allow override via container env.
    ; All STT_* vars from the container env are inherited by supervisord itself.
autostart=true
autorestart=true
startsecs=10
startretries=10
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth]
command=/usr/local/bin/node src/index.js
directory=/opt/auth
; Wait for STT to be ready before launching auth (auth healthchecks /ready against STT).
; start.sh handles the actual wait loop; supervisord just respawns if it dies.
environment=
    PORT="4000",
    NODE_ENV="production",
    DATABASE_URL="file:/var/lib/auth/prod.db",
    ; STT_ENDPOINT points at the in-process uvicorn — same container, localhost.
    STT_ENDPOINT="http://127.0.0.1:8000",
    ; JWT_RESET_SECRET / JWT_ACCESS_SECRET / JWT_REFRESH_SECRET inherited from container env.
autostart=false
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=false
autorestart=true
startsecs=1
startretries=5
stopsignal=QUIT
stopwaitsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
