# app/supervisor.conf
; STT-Endpoint all-in-one container — supervisord program definitions.
;
; Order of events:
;   1. supervisord starts (PID 1, but in all-in-one it is launched by start.sh
;      in the background and is NOT PID 1 — start.sh execs/wait's for it).
;   2. program:stt has autostart=true → uvicorn is up immediately.
;   3. start.sh waits for /ready, then issues supervisorctl start auth + nginx.
;   4. SIGTERM from the orchestrator reaches start.sh which forwards it to
;      supervisord; supervisord then tears down programs in stop-order.
;
; Log strategy (validated with supervisor 4.2.2 on Debian 11):
;   * [supervisord] logfile=/dev/null  → supervisord's own log goes nowhere.
;     We don't need it because program logs are the ones we care about.
;   * childlogdir=/var/log/supervisord → required real directory (set in
;     Dockerfile.allinone).
;   * stdout_logfile=/dev/stdout and stderr_logfile=/dev/stderr per program
;     → child output appears on the container's stdout/stderr, captured by
;     Docker/Easypanel as the program logs.

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
logfile_backups=0
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisord
user=root
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:stt]
command=/usr/local/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --log-level info
; cwd=/opt/stt (parent of the app/ package). Setting cwd to /opt/stt/app
; would add /opt/stt/app to sys.path, which makes ``import queue`` resolve to
; /opt/stt/app/queue.py (our app module) instead of the stdlib's queue —
; producing the "attempted relative import with no known parent package"
; ImportError we saw in the first deploy. As the parent, Python adds
; /opt/stt and finds app.queue via the package path.
directory=/opt/stt
environment=
    STT_MODELS_DIR="/var/lib/stt/models",
    PYTHONUNBUFFERED="1",
    PYTHONDONTWRITEBYTECODE="1",
    ; All STT_* env vars are inherited from the container env so callers can
    ; override STT_MODEL_SIZE / STT_DEVICE / STT_COMPUTE_TYPE etc. via
    ; `docker run -e STT_MODEL_SIZE=...` or via docker-compose.
autostart=true
autorestart=true
startsecs=10
startretries=10
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth]
command=/usr/local/bin/node src/index.js
directory=/opt/auth
; auth depends on STT /ready but the wait is performed by start.sh, so
; supervisord only tracks the auth process from then on (autostart=false).
environment=
    PORT="4000",
    NODE_ENV="production",
    DATABASE_URL="file:/var/lib/auth/prod.db",
    STT_ENDPOINT="http://127.0.0.1:8000",
    ; JWT_RESET_SECRET / JWT_ACCESS_SECRET / JWT_REFRESH_SECRET / CSRF_SECRET
    ; are inherited from the container environment. Without them set, auth
    ; refuses to start (see auth-service/src/lib/env.js).
autostart=false
autorestart=true
startsecs=3
startretries=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=false
autorestart=true
startsecs=1
startretries=5
stopsignal=QUIT
stopwaitsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
